# dot_config/systemd/user/change-wallpaper.service
[Unit]
Description=Change Wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/change-wallpaper
