# dot_config/systemd/user/change-wallpaper.timer
[Unit]
Description=Change Wallpaper Every Hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
